# godoc.service
[Service]
ExecStart=%h/go/bin/socket-activated-godoc
